# voxmd.service
[Unit]
Description=VOX-MD WhatsApp Bot
After=network.target

[Service]
Type=simple
User=<user>
Group=<group>
WorkingDirectory=/path/to/voxmd
ExecStart=/bin/bash /path/to/voxmd/start.sh
Restart=on-failure
RestartSec=10s
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=voxmd
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target